# docker/compare-pdf/rootfs/etc/periodic/15min/compare_cleanup
#!/usr/bin/env bash

set -eu

find "$COMPARE_OUTPUT_DIR" -mindepth 1 -maxdepth 1 -mmin +${COMPARE_OUTPUT_CLEANER_AGE} -exec rm -rf {} \; 2>&1 >> "$COMPARE_OUTPUT_DIR/cleanup.log"
