# docker/export/rootfs/etc/periodic/15min/export_cleanup
#!/usr/bin/env bash

set -eu

find "$EXPORT_OUTPUT_DIR" -mindepth 1 -maxdepth 1 -mmin +${EXPORT_OUTPUT_CLEANER_AGE} -exec rm -rf {} \; 2>&1 >> "$EXPORT_OUTPUT_DIR/cleanup.log"
